# branch/sprint-ietf-73/legacy/scripts/old-unused/cronjobs/crontab.crontab
MAILTO=logs
#0 23 * * * /home/mlee/RELEASE/gen_wg_chairs_aliases.pl
#0 2 * * * /home/mlee/RELEASE/wgchairs-generate.pl; /home/ietf/scripts/wgchairs-listupdate.sh; /usr/local/mailman/bin/list_lists > /home/ietf/LOGs/lists.txt; /home/mlee/RELEASE/pop_imported_mailing_list.pl; /home/mlee/RELEASE/wgchairs_by_area.pl; /home/mlee/RELEASE/bofchairs-generate.pl; /home/mlee/RELEASE/generate_iesg_only.pl;
#
# to be removed after migration
#
#40 */4 * * * /home/mirror/bin/mirror-text
#
#10 20 * * * /home/mirror/bin/mirror-spam
#10 18 * * * /home/mirror/bin/mirror-web
#
# wget-MAILING-LISTS.sh
#
0 2 * * *  /home/mirror/bin/wget-MAILING-LISTS.sh
#
# generate the web archives
#
#10 */4 * * * nice -n 10 /home/mirror/bin/mhonarc-archives.cron > /tmp/imailarc-archives.log 2>&1
#10 */4 * * * /home/mirror/bin/mhonarc-archives.cron > /tmp/imailarc-archives.log 2>&1
10 * * * * /home/mirror/bin/mhonarc-loop.sh >> /tmp/imailarc-archives.log 2>&1
# 40 21 * * * nice -n 20 /home/mirror/bin/mhonarc-spam.cron > /tmp/imailarc-spam.log 2>&1
#
# check for writable ids
#
33 1 * * * /home/mirror/bin/checkforwritable > /dev/null 2>&1
#
#
1 0 1 * * /home/mirror/bin/rotate-mailing-lists
1 0 1 * * /home/mirror/bin/rotate-secure-mailing-lists
#
# we are mhonarcing hourly now
# 0 22 * * *  /home/mirror/bin/mhonarc-archives.id
